# jibri/rootfs/etc/systemd/system/jibri-icewm.service.d/override.conf
[Service]
ExecStartPost=
ExecStartPost=/usr/bin/start-pulseaudio-x11