# annolens-api/rootfs/etc/services.d/nodejs/run
#!/usr/bin/with-contenv bash

exec node app.js
